# Mac/Software/IntelliJIDEALicenseServer/run-license-server.cron
@reboot /Users/bovenson/Git/notes/Mac/Software/IntelliJIDEALicenseServer/run-license-server.sh
