# Computer Science/Linux/Codes/cron-test-01.cron
* * * * * /home/public/Git/notes/Linux/Codes/cron-test-01.sh
