# Other/git-push.cron
0 23 * * * cd /home/bovenson/Git/notes ; bash git-push.sh
