# git-push.cron
0 23 * * * bash /home/bovenson/Git/notes/git-push.sh
